# .crontab
# minute hour day_of_month month day_of_week command

# reset nvim keys log
@reboot     /bin/echo > /home/ceigh/.nvim_keys.log

# clear trash
@reboot     /bin/rm -rf /home/ceigh/.local/share/Trash
